# roles/duckdns/templates/duckdns.timer
[Unit]
Description=Update DuckDNS every 10 minutes

[Timer]
OnBootSec=10
OnUnitActiveSec=600

[Install]
WantedBy=timers.target
